# build/scorecard_qaci01.conf
[program:scorecard]
environment=DJANGO_SETTINGS_MODULE='scorecard.settings.qaci01',PATH='/usr/local/virtualenvs/scorecard/',HOME='/home/wicqacip'
directory=/opt/scorecard
command = /usr/local/virtualenvs/scorecard/bin/uwsgi --master --chdir /opt/scorecard/ --module scorecard.wsgi:application --socket /opt/scorecard/run/scorecard.sock --chmod-socket=666 --processes 3 --virtualenv /usr/local/virtualenvs/scorecard/ --pythonpath /usr/local/virtualenvs/scorecard/ --buffer-size 4096 --env DJANGO_SETTINGS_MODULE='scorecard.settings.qaci01'

user = wicqacip
group = wicqacip
autostart=true
autorestart=true
stopsignal=QUIT
stderr_logfile = /var/log/scorecard/scorecard.log
stdout_logfile = /var/log/scorecard/scorecard.log
redirect_stderr = true



[program:scorecard-celery]
environment=DJANGO_SETTINGS_MODULE='scorecard.settings.qaci01',PATH='/usr/local/virtualenvs/scorecard/'
directory=/opt/scorecard/
command=/usr/local/virtualenvs/scorecard/bin/celery worker -n qaci01.scorecard  --app=scorecard -Q scorecard_queue --logfile=/var/log/scorecard/scorecard-celery.log -l warn -c 3
user=wicqacip
autostart=true
autorestart=true
redirect_stderr=True

[program:scorecard-beat]
environment=DJANGO_SETTINGS_MODULE='scorecard.settings.qaci01',PATH='/usr/local/virtualenvs/scorecard/'
directory=/opt/scorecard/
command=/usr/local/virtualenvs/scorecard/bin/celery beat --app=scorecard --logfile=/var/log/scorecard/scorecard-beat.log -l warn
user=wicqacip
autostart=true
autorestart=true
redirect_stderr=True
